# scripts/crontab
* * * * * /app/scripts/run_scripts.sh >> /var/log/cron.log 2>&1
